# src/build/linux/debian_jessie_i386-sysroot/debian/libnss-db/DEBIAN/postinst
#!/bin/sh

set -e

if [ "$1" = "configure" ]; then
	if [ -x /etc/init.d/nscd ]; then
		invoke-rc.d nscd restart
	fi
fi

# Automatically added by dh_makeshlibs
if [ "$1" = "configure" ]; then
	ldconfig
fi
# End automatically added section


exit 0
